# irfan /JobSystem_Kit/launchagents/com.irfan.jobapply.watch.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.irfan.jobapply.watch</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/zsh</string>
		<string>-lc</string>
		<string>bash "$HOME/jobs/bin/auto_apply_watch.sh"</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardErrorPath</key>
	<string>/Users/irfan/jobs/logs/auto_apply_watch.log</string>
	<key>StandardOutPath</key>
	<string>/Users/irfan/jobs/logs/auto_apply_watch.log</string>
	<key>StartInterval</key>
	<integer>600</integer>
</dict>
</plist>
