# irfan /JobSystem_Kit/launchagents/com.irfan.dailyreport.email.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
 "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0"><dict>
  <key>Label</key><string>com.irfan.dailyreport.email</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>bash "$HOME/jobs/bin/job_daily_email.sh"</string>
  </array>
  <key>StartCalendarInterval</key>
  <dict><key>Hour</key><integer>18</integer><key>Minute</key><integer>20</integer></dict>
  <key>RunAtLoad</key><true/>
  <key>StandardOutPath</key><string>/Users/irfan/jobs/logs/daily_report_email.log</string>
  <key>StandardErrorPath</key><string>/Users/irfan/jobs/logs/daily_report_email.log</string>
</dict></plist>
